# server/vibe-radar-likes.service
[Unit]
Description=Vibe Coding Radar likes API
After=network-online.target
Wants=network-online.target

[Service]
Type=simple
User=www-data
Group=www-data
WorkingDirectory=/opt/vibe-radar-likes
EnvironmentFile=/etc/vibe-radar-likes.env
Environment=PYTHONDONTWRITEBYTECODE=1
ExecStart=/usr/bin/python3 /opt/vibe-radar-likes/likes_service.py
Restart=on-failure
RestartSec=3
NoNewPrivileges=true
PrivateTmp=true
ProtectHome=true
ProtectSystem=strict

[Install]
WantedBy=multi-user.target
